# crontab
* * * * * /home/dipankar/.config/i3/scripts/battery.sh >> /home/dipankar/battery.log 2>&1
* * * * * echo "abcd" >> /home/dipankar/cron-dipankar.txt 2>&1

