# linux-i3/crontab
* * * * * /home/dipankar/.config/i3/scripts/battery.sh >> /var/log/auto-battery.log 2>&1

